# config/supervisor.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:app]
directory=/app/backend
command=python run.py

[program:npm]
directory=/app/frontend
command=npm run development